# wcd.service
[Unit]
Description=Wallpaper Changer Daemon

[Service]
ExecStart=/home/bruno/.local/bin/wcd

[Install]
WantedBy=default.target
